# addons/images/rsu_ping_fetch/crontab
*/5 * * * * /usr/local/bin/python3 /home/pingfetch/rsuStatusFetch.py
