# services/addons/images/iss_health_check/crontab
PYTHONPATH=/home
0 */6 * * * /usr/local/bin/python3 /home/iss_health_checker.py
