# services/addons/images/count_metric/crontab
PYTHONPATH=/home
0 7 * * * /usr/local/bin/python3 /home/daily_emailer.py
0 0 * * * /usr/local/bin/python3 /home/mongo_counter.py
