# services/addons/images/rsu_ping_fetch/crontab
PYTHONPATH=/home
*/1 * * * * /usr/local/bin/python3 /home/rsuStatusFetch.py
0 0 * * * /usr/local/bin/python3 /home/purger.py
