# services/addons/images/rsu_ping_fetch/crontab
PYTHONPATH=/home/pingfetch
*/1 * * * * /usr/local/bin/python3 /home/pingfetch/rsuStatusFetch.py
0 0 * * * /usr/local/bin/python3 /home/pingfetch/purger.py
